# dashboard/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:indexer]
command=/usr/local/bin/dfpn-indexer
    --rpc-url %(ENV_SOLANA_RPC_URL)s
    --bind 0.0.0.0:3030
    --index-path /data/indexes
    --content-registry %(ENV_CONTENT_REGISTRY)s
    --analysis-marketplace %(ENV_ANALYSIS_MARKETPLACE)s
    --model-registry %(ENV_MODEL_REGISTRY)s
    --worker-registry %(ENV_WORKER_REGISTRY)s
    --rewards %(ENV_REWARDS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
